# app/proxy/coder-proxy.service
[Unit]
Description=coder-proxy – reverse proxy for code-server Unix socket
After=network.target
Requires=code-server.service

[Service]
Type=simple
Restart=always
ExecStart=/opt/code-server/bin/coder-proxy \
    -socket /var/apps/coder/target/code-server.sock \
    -proxy-socket /var/apps/coder/target/coder-proxy.sock \
    -prefix /app/coder
User=1000
Group=1001

[Install]
WantedBy=multi-user.target
